# systemd/user/motivation.service
[Unit]
Description=Reminder notification

[Service]
Type=oneshot
ExecStart=/home/bala/.local/bin/reminder.sh
